# cron/cron.cron
* * * * * touch /var/www/RFHarvester/a.txt
1 * * * * touch /var/www/RFHarvester/b.txt
28 14 * * * (cd /var/www/RFHarvester && java -jar -Xms1024m -Xmx2048m -Dfile.encoding=UTF-8 bin/RFHarvesterV2test.jar --debug --harvest=portfolio)